# gateway/claire.gateway.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>claire.gateway</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/node</string>
        <string>/Users/sergio/sentientsergio/claire/gateway/dist/index.js</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/sergio/sentientsergio/claire/gateway</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>GATEWAY_PORT</key>
        <string>18789</string>
        <key>WORKSPACE_PATH</key>
        <string>/Users/sergio/sentientsergio/claire/workspace</string>
        <!-- ANTHROPIC_API_KEY should be set in your shell profile or here -->
    </dict>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/Users/sergio/Library/Logs/claire/gateway.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/sergio/Library/Logs/claire/gateway.error.log</string>
    
    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
