# gateway/claire.gateway.dev.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>claire.gateway.dev</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/node</string>
        <string>/Users/sergio/sentientsergio/claire/gateway/dist/index.js</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/sergio/sentientsergio/claire/gateway</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>NODE_ENV</key>
        <string>development</string>
    </dict>
    
    <key>RunAtLoad</key>
    <false/>
    
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/Users/sergio/Library/Logs/claire/gateway.dev.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/sergio/Library/Logs/claire/gateway.dev.error.log</string>
    
    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
